# deploy/wwm/systemd/mindchain-wwm-external-probe.timer
[Unit]
Description=Run MindChain WWM external probe every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=15s
Persistent=true
Unit=mindchain-wwm-external-probe.service

[Install]
WantedBy=timers.target
