# deploy/wwm/systemd/mindchain-wwm-witness@.service
[Unit]
Description=MindChain WWM public testnet witness %i
Documentation=https://wwm-status.mindchain.network/status.json
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=mindchain-wwm
Group=mindchain-wwm
EnvironmentFile=/etc/mindchain-wwm/witness-%i.env
ExecStart=/opt/mindchain-wwm/bin/mindchain-wwm-seed-launcher.sh
Restart=always
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
ReadWritePaths=/var/lib/mindchain-wwm-witness-%i

[Install]
WantedBy=multi-user.target
